# launchd/com.btcgenius.whale-stream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.btcgenius.whale-stream</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/tatianathevisionary/market-genius/src/whale_stream.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/tatianathevisionary/market-genius</string>

    <!-- Always-on websocket daemon: restart if it ever exits -->
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>

    <!-- Back off 10s between restarts if it's crash-looping -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/tatianathevisionary/market-genius/logs/whale_stream.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tatianathevisionary/market-genius/logs/whale_stream.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
